# rpi/systemd.service
[Unit]
Description=Tree-Project

[Service]
ExecStart=/home/pi/shula/rpi/init-tree-project.sh
Restart=always

[Install]
WantedBy=multi-user.target
